# raspberry-nfc-bridge/nfc-bridge.service
[Unit]
Description=NFC Bridge Server per Loyalty System
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/nfc-bridge
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=NFC_BRIDGE_PORT=3001

# Log output
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target